# src/freenas/etc/init.d/ix-aliases
#!/sbin/openrc-run
#
# $FreeBSD$
#
name="ix-aliases"

depend()
{
	provide ix-aliases
	before mountlate
}

start()
{
	local IFS="|"
	local f="bsdusr_username bsdusr_email"
	eval local $f
	local sf=$(var_to_sf $f)
	local user email
	cp /conf/base/etc/aliases /etc/aliases
	RO_FREENAS_CONFIG=$(ro_sqlite ${name} 2> /tmp/${name}.fail && rm /tmp/${name}.fail)
	trap 'rm -f ${RO_FREENAS_CONFIG}' EXIT
	${FREENAS_SQLITE_CMD} ${RO_FREENAS_CONFIG} "SELECT $sf FROM account_bsdusers WHERE bsdusr_email != ''" | \
	while eval read $f; do
		if grep -q "^${bsdusr_username}" /etc/aliases; then
			sed -i '' -E "s/^${bsdusr_username}:.*/${bsdusr_username}: ${bsdusr_email}/" /etc/aliases
		else
			echo "${bsdusr_username}: ${bsdusr_email}" >> /etc/aliases
		fi
	done
}
