# src/freenas/etc/init.d/ix-etc
#!/sbin/openrc-run
#
# $FreeBSD$
#

name="ix-etc"

depend()
{
	provide ix-etc
	need earlykld
	before fsck
}

start()
{
	LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/midclt call etc.generate_all > /dev/null
}
