# src/freenas/etc/init.d/ix-fstab
#!/sbin/openrc-run
#
# $FreeBSD$
#

name="ix-fstab"

depend()
{
	provide ix-fstab
	need earlykld
	before fsck
}

#
# Generate fstab right before mountlate.
#
generate_fstab_real()
{
	echo "fdescfs	/dev/fd	fdescfs rw	0 0"
}

generate_fstab_swap()
{
	local devname
	LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/midclt call disk.swaps_configure | LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/jq .[] | while read -r devname; do
		devname=$(eval echo $devname)
		echo "/dev/${devname}.eli	none	swap	sw	0	0"
	done;
}

start()
{
	freenas=$(LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/midclt call system.is_freenas)
	local _doumount=0
	RO_FREENAS_CONFIG=$(ro_sqlite ${name} 2> /tmp/${name}.fail && rm /tmp/${name}.fail)
	trap 'rm -f ${RO_FREENAS_CONFIG}' EXIT

	if [ `stat -f "%d:%i" /conf/base/etc/fstab 2>/dev/null` = `stat -f "%d:%i" /etc/fstab 2>/dev/null` ]; then
		# Hard link.  Let's remove /etc/fstab and copy /conf/base/etc/fstab to it
		rm /etc/fstab 2>/dev/null
	fi
	cp /conf/base/etc/fstab /etc/fstab 2>/dev/null
	generate_fstab_real >> /etc/fstab
	generate_fstab_swap >> /etc/fstab
	mount -uw /
}

