# src/freenas/etc/init.d/ix-localtime
#!/sbin/openrc-run
#
# $FreeBSD$
#

name="ix-localtime"

depend()
{
	provide ix-localtime
	before net devd
}


#
# Set up timezone
#
start()
{
	local IFS="|"
	local f="stg_timezone"
	eval local $f
	local sf=$(var_to_sf $f)
	RO_FREENAS_CONFIG=$(ro_sqlite ${name} 2> /tmp/${name}.fail && rm /tmp/${name}.fail)
	trap 'rm -f ${RO_FREENAS_CONFIG}' EXIT
	${FREENAS_SQLITE_CMD} ${RO_FREENAS_CONFIG} "SELECT $sf FROM system_settings ORDER BY -id LIMIT 1" | \
	while eval read $f; do
		timezone=$(echo ${stg_timezone} | sed -es=america-los_angeles=America/Los_Angeles=)
		if [ -z "${stg_timezone}" ]; then
			stg_timezone=America/Los_Angeles
		fi
		cp /usr/share/zoneinfo/${stg_timezone} /etc/localtime
		echo ${stg_timezone} > /var/db/zoneinfo
	done
}

