# nas_ports/freenas/py-middlewared/files/openrc-middlewared.in
#!/sbin/openrc-run
#
# $FreeBSD$
#

depend()
{
	before fsck earlykld
	provide middlewared
	keyword -shutdown
}

: ${middlewared_plugins_dirs=""}

# Set location of aditional plugin dirs
plugins_dirs_arg=""
for p in ${middlewared_plugins_dirs}; do
	plugins_dirs_arg="${plugins_dirs_arg} -p ${p}"
done

name="middlewared"
pidfile=/var/run/middlewared.pid
supervisor=supervise-daemon
command="/usr/local/bin/middlewared"
command_args="$plugins_dirs_arg --log-handler=file"
export LD_LIBRARY_PATH="/usr/local/lib"
export PATH="$PATH:/usr/local/sbin:/usr/local/bin"
export LC_ALL="en_US.UTF-8"

start_pre() {
	mount -uw /  # FIXME: why root is not rw by default?
	if ! mount | grep -q 'on /dev/fd (' 2> /dev/null; then
		mount -t fdescfs fdescfs /dev/fd
	fi
	/sbin/ifconfig lo0 127.0.0.1 up
	else
		supervise-daemon -p ${pidfile} --start ${command} ${command_args}
	fi
}

start_post() {
	LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/midclt -t 120 waitready
}
